# infrastructure/monitoring/protocol/infra/protocol-analyst.timer
[Unit]
Description=Run Protocol Night Analyst at 00:00 UTC

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
